# dev-sbc/sbc-status-leds/files/x4-n100-sbc-status-leds.crontab
*/3 * * * *   root   /usr/sbin/x4-n100-sbc-status-leds

