# dev-embedded/sbc_status_leds/files/sbc_status_leds.crontab
*/3 * * * *   root   /usr/local/sbin/sbc_status_leds

