# dev-embedded/sbc_status_leds/files/sbc_status_leds.crontab
*/3 * * * *   root   /usr/sbin/sbc_status_leds

